# provision/roles/restic/files/restic-nas-backup.timer
[Unit]
Description=Restic Weekly Backups (NAS)
[Timer]
OnCalendar=weekly
Persistent=true
[Install]
WantedBy=timers.target
