# provision/roles/restic/files/restic-b2-prune.timer
[Unit]
Description=Prune data from Restic B2 Backup
[Timer]
OnCalendar=monthly
Persistent=true
[Install]
WantedBy=timers.target
