# roles/restic/templates/restic-b2-backup.timer.j2
[Unit]
Description=Restic B2 Weekly Backups
[Timer]
OnCalendar=weekly
Persistent=true
[Install]
WantedBy=timers.target
