# provision/files/restic/restic-backup.timer
[Unit]
Description=Restic Weekly Backups
[Timer]
OnCalendar=weekly
Persistent=true
[Install]
WantedBy=timers.target
